# crontab
0 2 * * * /usr/local/bin/job.sh >> /var/log/cron.log 2>&1
#